# k8s/systemd/crio-removal.service
[Unit]
Description=CRI-O removal service

[Service]
Type=oneshot
ExecStart=/bin/sh -c "/usr/local/bin/crio-removal.sh"
RemainAfterExit=true