# k8s/systemd/crio-installer-flatcar.service
[Unit]
Description=CRI-O installer service

[Service]
Type=oneshot
ExecStart=/bin/sh -c "/opt/bin/crio-installer.sh"
RemainAfterExit=true