# k8s/systemd/crio-installer.service
[Unit]
Description=CRI-O installer service

[Service]
Type=oneshot
ExecStart=/bin/sh -c "/usr/local/bin/crio-installer.sh"
RemainAfterExit=true