# k8s/systemd/crio-removal-flatcar.service
[Unit]
Description=CRI-O removal service

[Service]
Type=oneshot
ExecStart=/bin/sh -c "/opt/bin/crio-removal.sh"
RemainAfterExit=true