# k8s/systemd/kubelet-config-helper.service
[Unit]
Description=Kubelet config service

[Service]
Type=oneshot
ExecStart=/bin/sh -c "/usr/local/bin/kubelet-config-helper.sh"
RemainAfterExit=true